# crontab
#Ten minutes after frappe framework move to production mode
*/10 * * * *      /home/production-mode.sh
