# crontab
*/10 * * * *      /home/production-mode.sh
